# airootfs/etc/xdg/autostart/Calamares.desktop
[Desktop Entry]
Type=Application
Exec=sudo -E calamares